# django-template/{{ cookiecutter.project_name }}/conf/supervisor/rqworker.conf
[program:rqworker]
directory={{ cookiecutter.prj_dr }}
command=python manage.py rqworker
autorestart=true
user=www-data
startsecs=5
